# devify.rules
#SUBSYSTEM=="hid", ACTION=="add|remove", \
#RUN+="/usr/bin/devify \"$env{HID_NAME}\" \"$env{ACTION}\"", GOTO="end_lookup"

# Input device
SUBSYSTEM=="input", ACTION=="add|remove", ENV{DEVLINKS}!="" \
ENV{DEVIFY_TYPE}="input_device", ENV{DEVIFY_DESC}="Input device", \
RUN+="/usr/bin/devify" GOTO="end_lookup"

# Removable drive
SUBSYSTEM=="block", ENV{DEVNAME}=="/dev/sd?", ACTION=="add|remove", \
ENV{DEVIFY_TYPE}="removable_drive", ENV{DEVIFY_DESC}="Removable drive", \
RUN+="/usr/bin/devify" GOTO="end_lookup"

# Android device
SUBSYSTEM=="usb", ENV{ID_VENDOR_FROM_DATABASE}=="Google Inc.", \
ENV{ID_MODEL_FROM_DATABASE}=="Nexus/Pixel Device*", ACTION=="add|remove" \
ENV{DEVIFY_TYPE}="android_device" , ENV{DEVIFY_DESC}="Android device", \
RUN+="/usr/bin/devify" GOTO="end_lookup"

# USB serial device, tested with Arduino and PIC K150 Programmer
SUBSYSTEM=="tty", ENV{DEVLINKS}=="/dev/serial*", ACTION=="add|remove", \
ENV{DEVIFY_TYPE}="serial_device", ENV{DEVIFY_DESC}="Serial device", \
RUN+="/usr/bin/devify" GOTO="end_lookup"

# Power adapter
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="?", ENV{POWER_SUPPLY_TYPE}=="Mains", ACTION=="change", \
ENV{DEVIFY_TYPE}="power_adapter", ENV{DEVIFY_DESC}="Power Adapter", \
RUN+="/usr/bin/devify" GOTO="end_lookup"

# # Sound cards
SUBSYSTEM=="sound", ENV{SOUND_INITIALIZED}=="1", ACTION=="change|remove", \
ENV{DEVIFY_TYPE}="sound_card", ENV{DEVIFY_DESC}="Audio device", \
RUN+="/usr/bin/devify" GOTO="end_lookup"

# Network adapters / usb thethering
ACTION=="add|remove", SUBSYSTEM=="net", \
ENV{DEVIFY_TYPE}="network_adapter", ENV{DEVIFY_DESC}="Network adapter", \
RUN+="/usr/bin/devify" GOTO="end_lookup"

# Try to avoid matching more than one rule per udev event
# Calling GOTO="end_lookup" in each rule
# This should always be at the end of this file!
LABEL="end_lookup"
