# devify.rules
SUBSYSTEM=="hid", ACTION=="remove"|| ACTION=="bind" , \
RUN+="/usr/local/bin/devify \"$env{HID_NAME}\" \"$env{ACTION}\""
